# config/home;user;.config/systemd;user/obex.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/lib/bluetooth/obexd -a --root=Downloads/Bluetooth
